# deploy/systemd/wh-preview.service
[Unit]
Description=WH Repeater diagnostic HDMI preview
After=wh-pc-gateway.service
Requires=wh-pc-gateway.service

[Service]
Type=simple
Environment=WH_PREVIEW_URL=udp://127.0.0.1:15000
Environment=WH_PREVIEW_CONNECTOR=DP-1
ExecStartPre=-/usr/bin/curl -fsS -X POST http://127.0.0.1:8080/api/preview/start
ExecStart=/usr/bin/mpv --no-terminal --really-quiet --idle=yes --keep-open=yes --force-window=yes --vo=drm --drm-connector=${WH_PREVIEW_CONNECTOR} --hwdec=vaapi-copy ${WH_PREVIEW_URL}
ExecStopPost=-/usr/bin/curl -fsS -X POST http://127.0.0.1:8080/api/preview/stop
Restart=on-failure
RestartSec=1s
TimeoutStopSec=5s
KillMode=mixed

[Install]
WantedBy=multi-user.target
